# debian/pipresents-rsync.cron.d
* * * * * root /etc/init.d/pipresents update >> /tmp/pipresents-update.log

